# modules/spraints/templates/etc/init/visage.conf.erb
description "Visage"
author "Lindsay Holmwood"

start on virtual-filesystems
stop on runlevel [06]

env PATH=/usr/local/bin:/usr/bin:/bin

env RAILS_ENV=production
env CONFIG_PATH=<%= @visage_config_path %>

setuid visage
setgid visage

console log

chdir /opt/visage/

script
  visage_root=`bundle show visage`
  exec ./bin/unicorn --config-file unicorn.rb -p 12004 --env production ${visage_root}/lib/visage-app/config.ru
end script
