# modules/spraints/templates/etc/rc.d/sprouter_config.erb
daemon=<%= @install_dir %>/bin/rails

. /etc/rc.d/rc.subr

rc_reload=NO
pexp="ruby: puma .*"

rc_pre() {
  ${rcexec} "cd <%= @install_dir %> && bin/rake assets:precompile && export RAILS_ENV=production && bin/rake db:migrate
}

rc_start() {
  ${rcexec} "cd <%= @install_dir %> && export RAILS_ENV=production && ${daemon} server --daemon --binding=0.0.0.0 --port=4433 ${daemon_flags}"
}

rc_cmd $1
