# modules/spraints/templates/etc/init.d/fowsr-ubuntu.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          fowsr
# Required-Start:    $network $remote_fs $syslog
# Required-Stop:     $network $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO

NAME="fowsr"
DESC="Weather Station Data Reader Service"
ROOT="/opt/fowsr/fowsr.src"
SCRIPT="$ROOT/fowsr-server.rb"
DAEMON="/usr/bin/ruby2.3"
PIDFILE="/run/fowsr.pid"
SOCK="/run/fowsr.socket"

HELPER="/usr/lib/fowsr/package-helper"

test -x $DAEMON || exit 0
test -f $SCRIPT || exit 0

. /lib/lsb/init-functions

# Override this variable by editing or creating /etc/default/fowsr.
DAEMON_OPTS="--fowsr $ROOT/fowsr --listen $SOCK"

if [ -f /etc/default/fowsr ]; then
    . /etc/default/fowsr
fi

case "$1" in
    start)
        log_daemon_msg "Starting $DESC" "$NAME"
        if start-stop-daemon --start --quiet --oknodo --pidfile $PIDFILE --name $NAME --startas $DAEMON -- $SCRIPT $DAEMON_OPTS; then
            log_end_msg 0
        else
            log_end_msg 1
        fi
        ;;

    stop)
        log_daemon_msg "Stopping $DESC" "$NAME"
        if start-stop-daemon --stop --quiet --oknodo --pidfile $PIDFILE --name $NAME; then
            log_end_msg 0
        else
            log_end_msg 1
        fi
        ;;

    reload|restart|force-reload)
        start-stop-daemon --stop --quiet --pidfile $PIDFILE --name $NAME --retry 5
        if start-stop-daemon --start --quiet --oknodo --pidfile $PIDFILE --name $NAME --startas $DAEMON -- $SCRIPT $DAEMON_OPTS; then
            log_end_msg 0
        else
            log_end_msg 1
        fi
        ;;

    status)
        status_of_proc -p $PIDFILE $DAEMON $NAME && exit 0 || exit $?
        ;;

    *)
        N=/etc/init.d/$NAME
        echo "Usage: $N {start|stop|restart|status|reload|force-reload}" >&2
        exit 1
        ;;
esac

exit 0
