# modules/spraints/templates/etc/rc.local.erb
#!/bin/sh

<% @zig_routes.each do |net| -%>
route -qn add -net <%= net %> <%= @zig_gw %><% end %>
<% @att_routes.each do |net| -%>
route -qn add -net <%= net %> <%= @att_gw %><% end %>
<% @hbb_routes.each do |net| -%>
route -qn add -net <%= net %> <%= @hbb_gw %><% end %>
